# conf/etc/init/fangmi-api.conf
start on [2345]
stop on [06]

script
    cd /var/www/fangmi-api/
    exec env/bin/uwsgi --ini /etc/fangmi-api/fangmi-api.ini
end script
